# bin/ocp-quicklab.service
[Unit]
Description=OCP quicklab

[Service]
ExecStart=/usr/local/bin/ocp-quicklab
Restart=on-failure
RestartSec=10
startLimitIntervalSec=60

Environment="AWS_SHARED_CREDENTIALS_FILE=/root/.aws/credentials"

StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target